# crontab
* * * * *;on
*/5 * * * *;off
*/5 * * * *;off
*/10 * * * *;on
